# compiler/crontab
* * * * * find /app/tmp/vch/* -mtime +1 -exec rm {} \;
